# supervisor/bash.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/bin/bash